# launchd/com.zdrowskit.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
Template only. Do not copy this file directly.

Run this from the repo root to generate a user-specific plist with your
checkout path, uv path, HOME, and log location:

    uv run python main.py daemon-install
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.zdrowskit.daemon</string>

    <key>ProgramArguments</key>
    <array>
        <string>__UV_PATH__</string>
        <string>run</string>
        <string>python</string>
        <string>__PROJECT_DIR__/src/daemon.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__PROJECT_DIR__</string>

    <key>KeepAlive</key>
    <true/>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/zdrowskit.daemon.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/zdrowskit.daemon.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>__PATH__</string>
        <key>HOME</key>
        <string>__HOME__</string>
    </dict>
</dict>
</plist>
